# variant_store/etc/service/ga4gh_server/run
#!/bin/sh
exec $VARIANT_STORE_DIR/.venv/bin/python $VARIANT_STORE_DIR/GA4GHAPI/ga4gh/runserver.py
