# docker/store/etc/service/ga4gh_server/run
#!/bin/sh
exec $STORE_DIR/GA4GHAPI/ga4gh/runserver.py
